# com.ofirblaicher.triage_report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ofirblaicher.triage_report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ofirblaicher/Documents/GitHub/bq_at_reports/run_report_auto.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ofirblaicher/Documents/GitHub/bq_at_reports/triage_report.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ofirblaicher/Documents/GitHub/bq_at_reports/triage_report.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
